# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/3_message/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

set +e
declare -a options

CONFIG_PATH=data/options.json


if bashio::config.true 'Send_Message3'; then
        etc/s6-overlay/s6-rc.d/3_message/run
fi
