# haos-sbfspot/rootfs/etc/s6-overlay/s6-rc.d/cron-service/run
#!/bin/execlineb -P
# Run cron in foreground so S6 can supervise it
cron -f
