# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/2_message/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

CONFIG_PATH=/data/options.json

MESSAGE3="$(bashio::config 'Message2')"
if bashio::config.true 'Send_Message2'; then
        /etc/s6-overlay/s6-rc.d/2_message/run
fi
