# haos-sbfspot/rootfs/etc/s6-overlay/s6-rc.d/sbfspot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SBFspot service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

CONFIG_PATH=/data/options.json


# ---- Symbolic Links ----
ln -sf /usr/bin/sbfspot/SBFspot /SBFspot
ln -sf /usr/bin/sbfspot/SBFspot.cfg /SBFspot.cfg
ln -sf /usr/bin/sbfspot/SBFspotUpload.cfg /SBFspotUpload.cfg

# SBFspot

exec /usr/bin/sbfspot/SBFspot -v -ad1 -am0 -ae0 -mqtt -finq
