# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/test/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

declare -a options

CONFIG_PATH=/data/options.json

MQTT_Host="$(bashio::config 'MQTT_Host' 'core-mosquitto')"
MQTT_User="$(bashio::config 'MQTT_User')"
MQTT_Pass="$(bashio::config 'MQTT_Pass')"
MQTT_Topic="$(bashio::config 'MQTT_Topic')"
MQTT_Data="$(bashio::config 'MQTT_Data')"

PLANTNAME="$(bashio::config 'Plantname')"
InvSerial="$(bashio::config 'InvSerial')"

PLANTNAME2="$(bashio::config 'Plantname2')"
InvSerial2="$(bashio::config 'InvSerial2')"

current_timestamp=$(date +'%Y-%m-%dT%H:%M:%S%z')


sleep 3

if bashio::config.true 'test_msg'; then 
        bashio::log.red "Sending Test Message..."

        bashio::log.green /usr/bin/mosquitto_pub -h "$MQTT_Host" -u "$MQTT_User" -P "$MQTT_Pass" -t homeassistant/sbfspot_"$PLANTNAME"/sbfspot_"$InvSerial2" -m "{\"PrgVersion\": \"3100\",\"Plantname\": \"SNA5000\",\"Timestamp\": \"$current_timestamp\",\"SunRise\": \"$current_timestamp\",\"SunSet\": \"$current_timestamp\",\"InvSerial\": \"2120108811\",\"InvName\": \"SN: 2120108811\",\"InvTime\": \"$current_timestamp\",\"InvStatus\": \"Ok\",\"InvSwVer\": \"031037R\",\"InvClass\": \"Solar Inverters\",\"InvType\": \"SB 5000TL-20\",\"InvTemperature\": \"0000\",\"InvGridRelay\": \"Information not available\",\"EToday\": \"6334\",\"ETotal\": \"75462350\",\"GridFreq\": \"2000\",\"PACTot\": \"0.000\",\"PAC1\": \"0000\",\"UAC1\": \"0000\",\"IAC1\": \"0000\",\"OperTm\": \"59117104\",\"FeedTm\": \"57784924\",\"PDCTot\": \"0000\",\"UDC1\": \"0000\",\"UDC2\": \"0000\",\"IDC1\": \"22\",\"IDC2\": \"454\",\"PDC1\": \"465400\",\"PDC2\": \"55000\",\"BTSignal\": \"68627\",\"InvWakeupTm\": \"$current_timestamp\",\"InvSleepTm\": \"$current_timestamp\"}" -d


        /usr/bin/mosquitto_pub -h "$MQTT_Host" -u "$MQTT_User" -P "$MQTT_Pass" -t homeassistant/sbfspot_"$PLANTNAME"/sbfspot_"$InvSerial2" -m "{\"PrgVersion\": \"3100\",\"Plantname\": \"SNA5000\",\"Timestamp\": \"$current_timestamp\",\"SunRise\": \"$current_timestamp\",\"SunSet\": \"$current_timestamp\",\"InvSerial\": \"2120108811\",\"InvName\": \"SN: 2120108811\",\"InvTime\": \"$current_timestamp\",\"InvStatus\": \"Ok\",\"InvSwVer\": \"031037R\",\"InvClass\": \"Solar Inverters\",\"InvType\": \"SB 5000TL-20\",\"InvTemperature\": \"0000\",\"InvGridRelay\": \"Information not available\",\"EToday\": \"6334\",\"ETotal\": \"75462350\",\"GridFreq\": \"2000\",\"PACTot\": \"0.000\",\"PAC1\": \"0000\",\"UAC1\": \"0000\",\"IAC1\": \"0000\",\"OperTm\": \"59117104\",\"FeedTm\": \"57784924\",\"PDCTot\": \"0000\",\"UDC1\": \"0000\",\"UDC2\": \"0000\",\"IDC1\": \"22\",\"IDC2\": \"454\",\"PDC1\": \"465400\",\"PDC2\": \"55000\",\"BTSignal\": \"68627\",\"InvWakeupTm\": \"$current_timestamp\",\"InvSleepTm\": \"$current_timestamp\"}" -d

fi
