# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/test/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

declare -a options

CONFIG_PATH=/data/options.json

current_timestamp=$(date +'%Y-%m-%dT%H:%M:%S%z')


sleep 8

if bashio::config.true 'test_msg'; then 
        bashio::log.red "Sending Test Message..."
        /usr/bin/mosquitto_pub -h "$(bashio::config 'MQTT_Host')" -u "$(bashio::config 'MQTT_Pass')" -t homeassistant/sbfspot_"$(bashio::config 'Plantname')"/sbfspot_"$(bashio::config 'InvSerial2')" -m "{\"PrgVersion\": \"3100\",\"Plantname\": \"SNA5000\",\"Timestamp\": \"$current_timestamp\",\"SunRise\": \"$current_timestamp\",\"SunSet\": \"$current_timestamp\",\"InvSerial\": \"2120108811\",\"InvName\": \"SN: 2120108811\",\"InvTime\": \"$current_timestamp\",\"InvStatus\": \"Ok\",\"InvSwVer\": \"031037R\",\"InvClass\": \"Solar Inverters\",\"InvType\": \"SB 5000TL-20\",\"InvTemperature\": \"0000\",\"InvGridRelay\": \"Information not available\",\"EToday\": \"6334\",\"ETotal\": \"75462350\",\"GridFreq\": \"2000\",\"PACTot\": \"0.000\",\"PAC1\": \"0000\",\"UAC1\": \"0000\",\"IAC1\": \"0000\",\"OperTm\": \"59117104\",\"FeedTm\": \"57784924\",\"PDCTot\": \"0000\",\"UDC1\": \"0000\",\"UDC2\": \"0000\",\"IDC1\": \"22\",\"IDC2\": \"454\",\"PDC1\": \"465400\",\"PDC2\": \"55000\",\"BTSignal\": \"68627\",\"InvWakeupTm\": \"$current_timestamp\",\"InvSleepTm\": \"$current_timestamp\"}" -d

fi
