# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/test/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================


Test_Msg="$(bashio::config 'Test_Msg')"

if bashio::config.true 'Test_Msg'; then
        /etc/s6-overlay/s6-rc.d/test/run
fi        