# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/test/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

if bashio::config.true 'test_msg'; then 
            /etc/s6-overlay/s6-rc.d/test/run
fi
