# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/message3/up
if bashio::config.true 'Send_Message3'; then
       /etc/s6-overlay/s6-rc.d/Message3/run
fi