# haos-sbf-mqtt-sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/message1/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

/usr/bin/mosquitto_pub -h "$MQTT_Host" -u "$MQTT_User" -P "$MQTT_Pass" -t homeassistant/sbfspot_"$PLANTNAME"/sbfspot_"$InvSerial" -m "$MESSAGE1" -d
