# haos-sbf-mqtt-sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/message1/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================


if bashio::config.true 'Send_Message1'; then
     ./message1/run
fi

