# noip-renewer/rootfs/etc/services.d/noip-renewer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SBFspot service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
set -x
CONFIG_PATH=/data/options.json

message=$(echo "PWD is current directory $(pwd)")
message=$(echo "The current time is $(date)")
# docker pull simaofsilva/noip-renewer:latest

EMAIL="$(jq --raw-output '.NOIP_Email' $CONFIG_PATH)"
PASSWORD="$(jq --raw-output '.NOIP_Password' $CONFIG_PATH)"
NOIP-TOTP-KEY="$(jq --raw-output '.NOIP_TOTP_KEY' $CONFIG_PATH)"
TRANSLATE="$(jq --raw-output '.Translate_Enabled' $CONFIG_PATH)"

# docker run --rm --env NO_IP_USERNAME=$(EMAIL) --env NO_IP_PASSWORD=$(PASSWORD) --env NO_IP_TOTP_KEY=$(NOIP-TOTP-KEY) --env TRANSLATE_ENABLED=$(TRANSLATE) simaofsilva/noip-renewer:latest 


