# haos-sbf-mqtt-sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/message2/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

bashio::log.cyan "Sending Message 2..."

sleep 10

/usr/bin/mosquitto_pub -h "$MQTT_Host" -u "$MQTT_User" -P "$MQTT_Pass" -t homeassistant/sbfspot_"$PLANTNAME"/sbfspot_"$InvSerial2" -m "$MESSAGE2" -d
