# haos-sbf-mqtt-sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/message2/up
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

if bashio::config.true 'Send_Message2'; then
        ./message2/run
fi
