# sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/mail/up
sensorsgen/rootfs/etc/s6-overlay/s6-rc.d/mail/run
